# M02/1r TRIMESTRE/Examen_Final_UF1/Hostal/todos_los_ficheros/Portera.cron
Me llamo Ivan
Higueras Menendez.

19-21 00 2,17 10 1 ism48086815 /home/users/inf/hism1/ism48086815/Documents/M02/Examen_Final_UF1/Hostal/AticoDeLuxe/Operario.sh

*/30 * * * 3,5 ism48086815 /home/users/inf/hism1/ism48086815/Documents/M02/Examen_Final_UF1/Hostal/AticoDeLuxe/Operario.sh
